# sigma/execution/TDL-EXE-000665.yml
title: High Risk Scheduled Task Disabled
id: TDL-EXE-000665
status: experimental
description: this rule detects activities associated with critical system components and processes, such
  as system restoration, security features, backup mechanisms, and updates. By detecting such events,
  administrators can ensure the system's health and security status remain intact.
references: []
author: TDL
date: '2024-01-01'
modified: '2026-04-27'
tags:
- attack.execution
- attack.t1053
- detection.auth
- detection.execution
logsource:
  service: security
  product: windows
detection:
  selection:
    EventID:
    - 47011
    - 4701
  condition: selection
falsepositives:
- Legitimate admin activity
- Authorized service accounts
level: high
custom:
  v4_id: '000665'
