# sigma/execution/TDL-EXE-000664.yml
title: High Risk Scheduled Task Deleted
id: TDL-EXE-000664
status: experimental
description: This detection rule identifies Event ID 4699 (a scheduled task was deleted) in conjunction
  with critical process names such as WindowsUpdate, BitLocker, Defender and . The focus is to monitor
  the deletion of scheduled tasks associated with key system operations, which could indicate malicious
  activity attempting to disrupt essential security and update processes. Alerts generated by this rule
  should be investigated promptly to ensure system integrity and security remain intact.”
references: []
author: TDL
date: '2024-01-01'
modified: '2026-04-27'
tags:
- attack.execution
- attack.t1053
- detection.auth
- detection.execution
logsource:
  service: security
  product: windows
detection:
  selection:
    EventID:
    - 46991
    - 4699
  condition: selection
falsepositives:
- Legitimate admin activity
- Authorized service accounts
level: high
custom:
  v4_id: '000664'
