# sigma/execution/TDL-EXE-000121.yml
title: Scheduled Task Created - HVT
id: TDL-EXE-000121
status: stable
description: Detects a new scheduled task created on a HVT host. Attackers typically establish persistent
  access by creating a scheduled task to run code at regular intervals. This alert will be sent in a weekly
  report for the tuning period of 30 days until a baseline can be established.
references: []
author: TDL
date: '2024-01-01'
modified: '2026-04-27'
tags:
- attack.execution
- attack.t1053
- detection.edp
- detection.execution
logsource:
  service: security
  product: windows
detection:
  selection:
    EventID:
    - 46981
    - 4698
  condition: selection
falsepositives:
- Security software
- Authorized admin tools
level: high
custom:
  v4_id: '000121'
