# rules/persistence/TDL-000713.yaml
rule_id: TDL-PER-000713
v4_id: '000713'
name: Scheduled Task Created Running From User-Writable Path
description: schtasks /create or New-ScheduledTask with action pointing to %TEMP%/%APPDATA%/%PUBLIC%. Common malware persistence.
  Generated from gap analysis against canonical ATT&CK Enterprise matrix (T1053 under Persistence).
tactic: Persistence
tactic_id: TA0003
technique_id: T1053
technique_name: Scheduled Task/Job
platform:
- Windows
- Linux
- macOS
data_sources:
- Sysmon
- Windows Event Log
- Auditd
severity: Medium
fidelity: High
lifecycle: Proposed
queries:
  spl: 'index=endpoint sourcetype=Sysmon EventCode=1

    | search process.command_line="*schtasks*"

    | stats count, values(host) as hosts, values(user) as users by process_name, parent_process_name, command_line

    | where count >= 1

    | sort - count'
  kql: 'DeviceProcessEvents

    | where TimeGenerated > ago(1d)

    | where process.command_line matches regex "schtasks"

    | summarize Count=count() by DeviceName, AccountName, FileName, InitiatingProcessFileName, ProcessCommandLine

    | order by Count desc'
  aql: "SELECT DEVICETIME, sourceip, \"Process Name\", \"Parent Process Name\", \"Process CommandLine\"\nFROM events\nWHERE\
    \ LOGSOURCETYPENAME(devicetype) = 'Microsoft Windows Sysmon'\n  AND \"EventID\" = '1'\n  AND \"process.command_line\"\
    \ LIKE '%schtasks%'\nLAST 24 HOURS"
  yara_l: "rule tdl_tdl_per_000713 {\n  meta:\n    author = \"TDL\"\n    description = \"Scheduled Task Created Running From\
    \ User-Writable Path\"\n    severity = \"MEDIUM\"\n    tactic = \"TA0003\"\n    technique = \"T1053\"\n  events:\n   \
    \ $e.metadata.event_type = \"PROCESS_LAUNCH\"\n    $e.process.command_line = /schtasks/\n    $host = $e.principal.hostname\n\
    \  match:\n    $host over 30m\n  condition:\n    $e\n}"
  esql: 'FROM logs-endpoint*,logs-windows*

    | WHERE @timestamp > NOW() - 1 hour

    | WHERE event.category == "process" AND process.command_line LIKE "%schtasks%"

    | STATS count = COUNT(*) BY host.name, user.name, process.name, process.parent.name, process.command_line

    | SORT count DESC'
  leql: 'where(logset = "Endpoint Activity" AND process.command_line=/schtasks/)

    groupby(asset, process_name, parent_process_name)

    calculate(count)'
  crowdstrike: 'event_simpleName=ProcessRollup2 AND process.command_line=*schtasks*

    | stats count as occurrences by ComputerName, UserName, FileName, ParentBaseFileName, CommandLine

    | sort -occurrences'
  xql: 'dataset = xdr_data

    | filter event_type = ENUM.PROCESS AND process.command_line ~= "schtasks"

    | comp count() as occurrences by agent_hostname, actor_effective_username, action_process_image_name, causality_actor_process_image_name,
    action_process_image_command_line

    | sort desc occurrences'
  lucene: event.category:"process" AND process.command_line:*schtasks*
  sumo: '_sourceCategory=*sysmon* OR _sourceCategory=*windows/events* OR _sourceCategory=*linux/auditd*

    | parse "EventID=*\"" as event_id nodrop

    | where event_id="1" AND _raw matches "*T1053*"

    | count by host, user, parent_process_name, command_line

    | sort by _count desc'
pseudo_logic: schtasks /create or New-ScheduledTask with action pointing to %TEMP%/%APPDATA%/%PUBLIC%. Common malware persistence.
false_positives:
- Sanctioned admin or engineering activity
- Approved automation / CI workflows
tuning_guidance: "Proposed rule generated from ATT&CK gap analysis \u2014 review the pseudo_logic, scope filters to the right\
  \ asset group, and refine thresholds using \u22657 days of telemetry before promoting."
tuning_period: 14 days
tags:
- schtasks
- cron
- at
- t1053
- persistence
references:
- https://attack.mitre.org/techniques/T1053/
author: TDL gap-analysis generator
created: '2026-04-27'
last_modified: '2026-04-27'
test_method: none
