# rules/persistence/TDL-000229.yaml
rule_id: TDL-PER-000229
v4_id: 000229
name: Startup File Modification - OSX
description: Detects when files related to configuring startup processes are modified in OSX. Attackers can create malicious
  plist files and place them in the LaunchAgent and LaunchDaemon directories, or modify the rc.common or launchd.conf files
  to configure the malicious process to run on the machines startup process to establish persistence.
tactic: Persistence
tactic_id: TA0003
technique_id: T1037
technique_name: Boot or Logon Initialization Scripts
platform:
- macOS
data_sources:
- MAC Operating System
- 'Event IDs: 4663'
severity: High
fidelity: High
lifecycle: Proposed
risk_score: 75
queries:
  spl: |
    index=endpoint sourcetype=Sysmon EventCode=1
    | search event_id IN ("4663")
    | stats count, values(host) as hosts, values(user) as users by process_name, parent_process_name, command_line
    | where count >= 1
    | sort - count
  kql: |
    DeviceProcessEvents
    | where TimeGenerated > ago(1d)
    | where event_id in ("4663")
    | summarize Count=count() by DeviceName, AccountName, FileName, InitiatingProcessFileName, ProcessCommandLine
    | order by Count desc
  aql: |
    SELECT DEVICETIME, sourceip, "Process Name", "Parent Process Name", "Process CommandLine"
    FROM events
    WHERE LOGSOURCETYPENAME(devicetype) = 'Microsoft Windows Sysmon'
      AND "EventID" = '1'
      AND "event_id" IN ('4663')
    LAST 24 HOURS
  yara_l: |
    rule tdl_tdl_per_000229 {
      meta:
        author = "TDL"
        description = "Startup File Modification - OSX"
        severity = "HIGH"
        tactic = "TA0003"
        technique = "T1037"
      events:
        $e.metadata.event_type = "PROCESS_LAUNCH"
        $e.event_id = /(4663)/
        $host = $e.principal.hostname
      match:
        $host over 30m
      condition:
        $e
    }
  esql: |
    FROM logs-endpoint*,logs-windows*
    | WHERE @timestamp > NOW() - 1 hour
    | WHERE event.category == "process" AND event_id IN ("4663")
    | STATS count = COUNT(*) BY host.name, user.name, process.name, process.parent.name, process.command_line
    | SORT count DESC
  leql: |
    where(logset = "Endpoint Activity" AND event_id IN [4663])
    groupby(asset, process_name, parent_process_name)
    calculate(count)
  crowdstrike: |
    event_simpleName=ProcessRollup2 AND event_id IN (4663)
    | stats count as occurrences by ComputerName, UserName, FileName, ParentBaseFileName, CommandLine
    | sort -occurrences
  xql: |
    dataset = xdr_data
    | filter event_type = ENUM.PROCESS AND event_id in ("4663")
    | comp count() as occurrences by agent_hostname, actor_effective_username, action_process_image_name, causality_actor_process_image_name, action_process_image_command_line
    | sort desc occurrences
  lucene: |
    event.category:"process" AND event_id:(4663)
  sumo: |
    _sourceCategory=*sysmon* OR _sourceCategory=*endpoint*
    | parse "EventID=*\"" as event_id
    | where event_id="1" AND event_id IN ("4663")
    | count by host, user, parent_process_name, command_line
    | sort by _count desc
false_positives:
- Security software
- Authorized admin tools
tuning_guidance: 'Tune per environment baseline. Tuning period: 7 days.'
tags:
- edp
- persistence
author: TDL
created: '2024-01-01'
last_modified: '2026-04-28'
test_method: none
pseudo_logic: '1 Windows Modification Event

  Where the Windows Event ID is in [567, 4663]

  Where the File Name is in [Library/LaunchDaemons/, Library/LaunchAgents/, etc/rc.common, etc/launchd.conf]

  Where the Process Name is NOT in [installd, shove, mdworker, mds, backup]'
suppression:
  window: 8 hours
tuning_period: 7 days
triage_steps:
- 'Confirm the alert metadata: timestamp, host, user, and the persistence artifact (registry path / scheduled task / service
  / config).'
- Inspect the file path, file hash, parent process, and any prior network/email delivery — verify the persistence target points
  to an unusual location or unsigned binary.
- Determine when the artifact was created and which process/account performed the write — pivot on that process for the original
  attack chain.
- 'Look for prior child-process activity from the persistence target: did anything actually launch from it?'
- 'Decide TP / FP: FP for sanctioned management software (RMM, vendor agents, custom installers) — verify against the asset
  register.'
- 'If TP: remove the persistence artifact, isolate the host, hunt for the same artifact across the fleet, and rotate any credentials
  touched by the persistent process.'
