# rules/execution/TDL-000121.yaml
rule_id: TDL-EXE-000121
v4_id: '000121'
name: Scheduled Task Created - HVT
description: Detects a new scheduled task created on a HVT host. Attackers typically establish persistent access by creating
  a scheduled task to run code at regular intervals. This alert will be sent in a weekly report for the tuning period of 30
  days until a baseline can be established.
tactic: Execution
tactic_id: TA0002
technique_id: T1053
technique_name: Scheduled Task/Job
platform:
- Windows
data_sources:
- Windows Operating System
- 'Event IDs: 4698'
severity: High
fidelity: High
lifecycle: Deployed
risk_score: 75
queries:
  spl: |
    index=endpoint sourcetype=Sysmon EventCode=1
    | search event_id IN ("4698")
    | stats count, values(host) as hosts, values(user) as users by process_name, parent_process_name, command_line
    | where count >= 1
    | sort - count
  kql: |
    DeviceProcessEvents
    | where TimeGenerated > ago(1d)
    | where event_id in ("4698")
    | summarize Count=count() by DeviceName, AccountName, FileName, InitiatingProcessFileName, ProcessCommandLine
    | order by Count desc
  aql: |
    SELECT DEVICETIME, sourceip, "Process Name", "Parent Process Name", "Process CommandLine"
    FROM events
    WHERE LOGSOURCETYPENAME(devicetype) = 'Microsoft Windows Sysmon'
      AND "EventID" = '1'
      AND "event_id" IN ('4698')
    LAST 24 HOURS
  yara_l: |
    rule tdl_tdl_exe_000121 {
      meta:
        author = "TDL"
        description = "Scheduled Task Created - HVT"
        severity = "HIGH"
        tactic = "TA0002"
        technique = "T1053"
      events:
        $e.metadata.event_type = "PROCESS_LAUNCH"
        $e.event_id = /(4698)/
        $host = $e.principal.hostname
      match:
        $host over 30m
      condition:
        $e
    }
  esql: |
    FROM logs-endpoint*,logs-windows*
    | WHERE @timestamp > NOW() - 1 hour
    | WHERE event.category == "process" AND event_id IN ("4698")
    | STATS count = COUNT(*) BY host.name, user.name, process.name, process.parent.name, process.command_line
    | SORT count DESC
  leql: |
    where(logset = "Endpoint Activity" AND event_id IN [4698])
    groupby(asset, process_name, parent_process_name)
    calculate(count)
  crowdstrike: |
    event_simpleName=ProcessRollup2 AND event_id IN (4698)
    | stats count as occurrences by ComputerName, UserName, FileName, ParentBaseFileName, CommandLine
    | sort -occurrences
  xql: |
    dataset = xdr_data
    | filter event_type = ENUM.PROCESS AND event_id in ("4698")
    | comp count() as occurrences by agent_hostname, actor_effective_username, action_process_image_name, causality_actor_process_image_name, action_process_image_command_line
    | sort desc occurrences
  lucene: |
    event.category:"process" AND event_id:(4698)
  sumo: |
    _sourceCategory=*sysmon* OR _sourceCategory=*endpoint*
    | parse "EventID=*\"" as event_id
    | where event_id="1" AND event_id IN ("4698")
    | count by host, user, parent_process_name, command_line
    | sort by _count desc
false_positives:
- Security software
- Authorized admin tools
tuning_guidance: 'Tune per environment baseline. Tuning period: 7 days.'
tags:
- edp
- execution
author: TDL
created: '2024-01-01'
last_modified: '2026-04-28'
test_method: historical
pseudo_logic: '1 Windows Scheduled Task Created Event

  Where the Event ID is in [602, 4698]

  Where the Source Username does NOT end with "$"

  Where the Destination Hostname is on the High Value Target Host List'
suppression:
  window: 8 hours
  fields:
  - Source Username
tuning_period: 7 days
triage_steps:
- 'Confirm the alert metadata: timestamp, host, user, and the executing process.'
- Inspect the parent → child process tree, command line arguments, and code-signing status — pay attention to short-lived
  processes, encoded arguments, and unusual parents (Office, browser, services).
- Pull EDR process tree for the alert PID; check the parent / grandparent and any child processes spawned within 5 minutes.
- 'Compare against the host''s baseline: is the binary / interpreter normally invoked here, by this user, at this time?'
- 'Decide TP / FP: TP if the execution chain has no business justification; FP if it traces to a packaging script, software
  install, or scheduled job.'
- 'If TP: contain the host, kill the offending process, snapshot memory if possible, and escalate to IR with the process tree
  attached.'
