# rules/execution/TDL-000121.yaml
rule_id: TDL-EXE-000121
v4_id: '000121'
name: Scheduled Task Created - HVT
description: Detects a new scheduled task created on a HVT host. Attackers typically establish persistent access by creating
  a scheduled task to run code at regular intervals. This alert will be sent in a weekly report for the tuning period of 30
  days until a baseline can be established.
tactic: Execution
tactic_id: TA0002
technique_id: T1053
technique_name: Scheduled Task/Job
platform:
- Windows
data_sources:
- Windows Operating System
- 'Event IDs: 4698'
severity: High
fidelity: High
lifecycle: Deployed
risk_score: 75
queries:
  spl: 'index=windows

    | where (EventCode=4698 OR EventCode=4698)

    | table _time, host, user, ProcessName, CommandLine'
  kql: 'SecurityEvent

    | where TimeGenerated >= ago(1d)

    | where EventID in (4698, 4698)

    | project TimeGenerated, Computer, Account, IpAddress, EventID, Activity'
  aql: "SELECT LOGSOURCENAME(logsourceid) AS log_source, username, hostname, QIDNAME(qid) AS event_name, eventcount\nFROM\
    \ events\nWHERE eventid = 4698\n  AND INOFFSET(starttime, LAST 300 SECONDS)\nGROUP BY username, hostname, QIDNAME(qid)\n\
    \nORDER BY starttime DESC"
  yara_l: "rule zs_exe_000121 {\n  meta:\n    author = \"TDL\"\n    description = \"Scheduled Task Created - HVT\"\n    severity\
    \ = \"HIGH\"\n    tactic = \"TA0002\"\n    technique = \"T1053\"\n\n  events:\n    $event.metadata.event_type = \"GENERIC_EVENT\"\
    \n    $event.metadata.product_event_type = \"4698\"\n    $event.principal.hostname = $host\n    $event.principal.user.userid\
    \ = $user\n\n  match:\n    $host over 5m\n\n  condition:\n    #event >= 1\n\n  outcome:\n    $risk_score = max(50)\n}"
  esql: 'FROM logs-windows.sysmon_operational*,logs-system.security*

    | WHERE @timestamp > NOW() - 1 day

    | WHERE event.code IN (''4698'')

    | STATS event_count = COUNT(*) BY host.name, user.name, event.code, @timestamp = BUCKET(@timestamp, 5m)

    | WHERE event_count >= 1

    | SORT event_count DESC'
  leql: 'where(event_id IN [4698])

    groupby(hostname, username)

    calculate(count AS event_count)'
  crowdstrike: 'event_simpleName=SecurityEvent

    | where EventCode=4698

    | table ComputerName, UserName, CommandLine, ContextTimeStamp'
  xql: 'dataset = xdr_data

    | filter event_type = "PROCESS" and action_evtlog_event_id in (4698)

    | filter _time > to_timestamp(now() - 86400, "epoch")

    | fields hostname, actor_effective_username, action_process_image_name, action_process_command_line, _time

    | sort desc _time'
  lucene: event_id:(4698)
false_positives:
- Security software
- Authorized admin tools
tuning_guidance: 'Tune per environment baseline. Tuning period: 7 days.'
tags:
- edp
- execution
author: TDL
created: '2024-01-01'
last_modified: '2026-04-27'
test_method: historical
pseudo_logic: '1 Windows Scheduled Task Created Event

  Where the Event ID is in [602, 4698]

  Where the Source Username does NOT end with "$"

  Where the Destination Hostname is on the High Value Target Host List'
suppression:
  window: 8 hours
  fields:
  - Source Username
tuning_period: 7 days
