# rules/execution/TDL-000665.yaml
rule_id: TDL-EXE-000665
v4_id: '000665'
name: High Risk Scheduled Task Disabled
description: this rule detects activities associated with critical system components and processes, such as system restoration,
  security features, backup mechanisms, and updates. By detecting such events, administrators can ensure the system's health
  and security status remain intact.
tactic: Execution
tactic_id: TA0002
technique_id: T1053
technique_name: Scheduled Task or Job
platform:
- Windows
data_sources:
- Microsoft Windows
- 'Event IDs: 4701'
severity: High
fidelity: High
lifecycle: Proposed
risk_score: 75
queries:
  spl: |-
    index=windows
    | where (EventCode=4701 OR EventCode=4701)
    | table _time, host, Account_Name, IpAddress, LogonType
  kql: |-
    SecurityEvent
    | where TimeGenerated >= ago(1d)
    | where EventID in (4701, 4701)
    | project TimeGenerated, Computer, Account, IpAddress, EventID, Activity
  aql: |-
    SELECT LOGSOURCENAME(logsourceid) AS log_source, username, hostname, QIDNAME(qid) AS event_name, eventcount
    FROM events
    WHERE eventid = 4701
      AND INOFFSET(starttime, LAST 300 SECONDS)
    GROUP BY username, hostname, QIDNAME(qid)

    ORDER BY starttime DESC
  yara_l: |-
    rule zs_exe_000665 {
      meta:
        author = "TDL"
        description = "High Risk Scheduled Task Disabled"
        severity = "HIGH"
        tactic = "TA0002"
        technique = "T1053"

      events:
        $event.metadata.event_type = "GENERIC_EVENT"
        $event.metadata.product_event_type = "4701"
        $event.principal.hostname = $host
        $event.principal.user.userid = $user

      match:
        $host over 5m

      condition:
        #event >= 1

      outcome:
        $risk_score = max(50)
    }
  esql: |-
    FROM logs-windows.sysmon_operational*,logs-system.security*
    | WHERE @timestamp > NOW() - 1 day
    | WHERE event.code IN ('4701')
    | STATS event_count = COUNT(*) BY host.name, user.name, event.code, @timestamp = BUCKET(@timestamp, 5m)
    | WHERE event_count >= 1
    | SORT event_count DESC
  leql: |-
    where(event_id IN [4701])
    groupby(hostname, username)
    calculate(count AS event_count)
  crowdstrike: |-
    event_simpleName=SecurityEvent
    | where EventCode=4701
    | table ComputerName, UserName, CommandLine, ContextTimeStamp
  xql: |-
    dataset = xdr_data
    | filter event_type = "PROCESS" and action_evtlog_event_id in (4701)
    | filter _time > to_timestamp(now() - 86400, "epoch")
    | fields hostname, actor_effective_username, action_process_image_name, action_process_command_line, _time
    | sort desc _time
  lucene: event_id:(4701)
false_positives:
- Legitimate admin activity
- Authorized service accounts
tuning_guidance: Tune per environment baseline.
tags:
- auth
- execution
author: TDL
created: '2024-01-01'
last_modified: '2026-04-27'
test_method: none
pseudo_logic: "1 Windows Security Event\nWhere Event ID 4701\nWhere Taskname contains either \n\"\"*\\\\Windows\\\\SystemRestore\\\
  \\SR*\"\" \n\"\"*\\\\Windows\\\\Windows\\\\Defender*\"\" \n\"\"*\\\\Windows\\\\BitLocker*\"\" \n\"\"*\\\\Windows\\\\WindowsBackup*\"\
  \" \n\"\"*\\\\Windows\\\\WindowsUpdate*\"\" \n\"\"*\\\\Windows\\\\UpdateOrchestrator\\\\Schedule*\"\" \n\"\"*\\\\Windows\\\
  \\ExploitGuard*\"\""
