# rules/privilege-escalation/TDL-000719.yaml
rule_id: TDL-PE-000719
v4_id: 000719
name: New Windows Service Created With Binary in User Path
description: sc create / New-Service / direct registry create under HKLM\SYSTEM\CurrentControlSet\Services with ImagePath
  in %TEMP%/%APPDATA%. Service runs as SYSTEM at next start. Generated from gap analysis against canonical ATT&CK Enterprise
  matrix (T1543 under Privilege Escalation).
tactic: Privilege Escalation
tactic_id: TA0004
technique_id: T1543
technique_name: Create or Modify System Process
platform:
- Windows
- macOS
- Linux
data_sources:
- Sysmon
- EDR
- Auditd
severity: High
fidelity: High
lifecycle: Proposed
risk_score: 75
queries:
  spl: 'index=endpoint sourcetype=Sysmon EventCode=1

    | search process.command_line="*sc-create*"

    | stats count, values(host) as hosts, values(user) as users by process_name, parent_process_name, command_line

    | where count >= 1

    | sort - count'
  kql: 'DeviceProcessEvents

    | where TimeGenerated > ago(1d)

    | where process.command_line matches regex "sc-create"

    | summarize Count=count() by DeviceName, AccountName, FileName, InitiatingProcessFileName, ProcessCommandLine

    | order by Count desc'
  aql: "SELECT DEVICETIME, sourceip, \"Process Name\", \"Parent Process Name\", \"Process CommandLine\"\nFROM events\nWHERE\
    \ LOGSOURCETYPENAME(devicetype) = 'Microsoft Windows Sysmon'\n  AND \"EventID\" = '1'\n  AND \"process.command_line\"\
    \ LIKE '%sc-create%'\nLAST 24 HOURS"
  yara_l: "rule tdl_tdl_pe_000719 {\n  meta:\n    author = \"TDL\"\n    description = \"New Windows Service Created With Binary\
    \ in User Path\"\n    severity = \"HIGH\"\n    tactic = \"TA0004\"\n    technique = \"T1543\"\n  events:\n    $e.metadata.event_type\
    \ = \"PROCESS_LAUNCH\"\n    $e.process.command_line = /sc-create/\n    $host = $e.principal.hostname\n  match:\n    $host\
    \ over 30m\n  condition:\n    $e\n}"
  esql: 'FROM logs-endpoint*,logs-windows*

    | WHERE @timestamp > NOW() - 1 hour

    | WHERE event.category == "process" AND process.command_line LIKE "%sc-create%"

    | STATS count = COUNT(*) BY host.name, user.name, process.name, process.parent.name, process.command_line

    | SORT count DESC'
  leql: 'where(logset = "Endpoint Activity" AND process.command_line=/sc-create/)

    groupby(asset, process_name, parent_process_name)

    calculate(count)'
  crowdstrike: 'event_simpleName=ProcessRollup2 AND process.command_line=*sc-create*

    | stats count as occurrences by ComputerName, UserName, FileName, ParentBaseFileName, CommandLine

    | sort -occurrences'
  xql: 'dataset = xdr_data

    | filter event_type = ENUM.PROCESS AND process.command_line ~= "sc-create"

    | comp count() as occurrences by agent_hostname, actor_effective_username, action_process_image_name, causality_actor_process_image_name,
    action_process_image_command_line

    | sort desc occurrences'
  lucene: event.category:"process" AND process.command_line:*sc-create*
  sumo: '_sourceCategory=*sysmon* OR _sourceCategory=*endpoint/edr* OR _sourceCategory=*linux/auditd*

    | parse "EventID=*\"" as event_id nodrop

    | where event_id IN ("12","13","14") AND _raw matches "*T1543*"

    | count by host, user, registry_key, registry_value, image

    | sort by _count desc'
pseudo_logic: sc create / New-Service / direct registry create under HKLM\SYSTEM\CurrentControlSet\Services with ImagePath
  in %TEMP%/%APPDATA%. Service runs as SYSTEM at next start.
false_positives:
- Sanctioned admin or engineering activity
- Approved automation / CI workflows
tuning_guidance: "Proposed rule generated from ATT&CK gap analysis \u2014 review the pseudo_logic, scope filters to the right\
  \ asset group, and refine thresholds using \u22657 days of telemetry before promoting."
tuning_period: 14 days
tags:
- service
- launchd
- systemd
- t1543
- privilege-escalation
references:
- https://attack.mitre.org/techniques/T1543/
author: TDL gap-analysis generator
created: '2026-04-27'
last_modified: '2026-04-27'
test_method: none
triage_steps:
- 'Confirm the alert metadata: timestamp, host, source user, and (if visible) the resulting privileged context.'
- "Inspect the parent \u2192 child process tree, command line arguments, and code-signing status \u2014 verify the privilege\
  \ gain was not part of an approved elevation flow (UAC consent, sudo, PAM ticket)."
- 'Pull preceding 1h of activity for the source user: was there earlier credential-access / discovery / exploit telemetry?'
- 'Validate token / session changes: new SeImpersonate-style privileges, new IAM role assumption, or a SYSTEM-level child
  process from a user-level parent.'
- 'Decide TP / FP: TP if there''s no change-management correlation and no parent admin process; FP if it''s a documented elevation
  by IT.'
- 'If TP: contain the host, revoke the elevated session, sweep for other escalations using the same primitive, and notify
  identity / IR teams.'
